# packages/agenc-worker/templates/launchd/ai.tetsuo.agenc-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  agenc-worker one-shot tick every 5 minutes (launchd, macOS).

  Install:
    npm install -g @tetsuo-ai/agenc-worker@REVIEWED_EXACT_VERSION
    command -v agenc-worker  (update ProgramArguments below if the path differs)
    cp ai.tetsuo.agenc-worker.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/ai.tetsuo.agenc-worker.plist

  Configure via ~/.config/agenc-worker/config.json (or add an
  EnvironmentVariables dict). The wallet MUST be a low-funded hot wallet -
  it is the worker's only spend authority.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.tetsuo.agenc-worker</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Deliberately no npx/network install at service start. -->
    <string>/usr/local/bin/agenc-worker</string>
    <string>once</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/private/tmp</string>
  <key>Umask</key>
  <integer>63</integer>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <!-- Add StandardOutPath/StandardErrorPath only under a private directory;
       predictable shared-/tmp log paths are intentionally not shipped. -->
</dict>
</plist>
