# packages/agenc-worker/templates/systemd/agenc-worker.timer
# Run `agenc-worker once` every 5 minutes (see agenc-worker.service for the
# install steps).
[Unit]
Description=Run the AgenC marketplace worker tick every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
# Ticks are cheap when there is nothing claimable; randomize a little so a
# fleet of workers doesn't sweep the RPC in lockstep.
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
